# systemd/slushy-trade-mcp.service
[Unit]
Description=slushy.trade MCP server (supporter-gated, fronts HyPaper)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debian
WorkingDirectory=/opt/slushy-trade-mcp
EnvironmentFile=/opt/slushy-trade-mcp/.env
# npm cache/config lookups during npx — keep a writable HOME under /opt.
Environment=HOME=/opt/slushy-trade-mcp
ExecStart=/usr/bin/npx tsx src/index.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slushy-trade-mcp

# Hardening. OAuth client registrations are now persisted under
# /var/lib/slushy-trade-mcp/clients.json (systemd auto-creates the dir
# with the service user as owner). Everything else stays in-memory.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/slushy-trade-mcp
StateDirectory=slushy-trade-mcp
StateDirectoryMode=0700
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
